# websockify.conf
[program:websockify]
command=/system/websockify.sh
autorestart=true
